# EthereumCustomScripts/Service/ipfsservice.service
[Unit]
Description=IPFS Service

[Service]
ExecStart=ipfs daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target